# job-monitor.service
[Unit]
Description=Job Monitor Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/job-monitor-bot
ExecStart=/usr/bin/python3 main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
